# systemd-user-services/systemd/user/media-stack-backup.service
[Unit]
Description=Backup media stack configs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/backup-media-stack.sh
